# yusuke-notify.service
[Unit]
Description=Pacman Update Notifier

[Service]
Type=simple
ExecStart=/usr/bin/yusuke

[Install]
WantedBy=default.target
